# rkhunter.cron
#!/bin/sh
output=`/usr/sbin/rkhunter -c --quiet --nocolors`
[ "$output" != "" ] && echo $output | /bin/mail -s 'rkhunter Daily Run' root
