# rkhunter.cron
#!/bin/sh
umask 027; /usr/sbin/rkhunter --update >/dev/null
/usr/sbin/rkhunter -c --report-warnings-only --cronjob
