# rkhunter.cron
#!/bin/sh
/usr/sbin/rkhunter --cronjob | /bin/mail -s 'rkhunter Daily Run' root
