# rkhunter.cron
#!/bin/sh
/usr/sbin/rkhunter -c --report-warnings-only --cronjob
